# runsvdir-start.service
[Unit]
Description=Runit Process Supervisor

[Service]
Type=simple
ExecStart=/usr/sbin/runsvdir-start
Restart=always
RestartSec=500ms
